# event-detect-service/event-detector.service
[Unit]
Description=Event Detection Service - YOLO Object Detection for Image Monitoring
After=multi-user.target network.target
Wants=multi-user.target

[Service]
Type=simple
# Update these paths according to your setup
ExecStart=/usr/bin/python3 /home/pi/Documents/watchmen/event-detect-service/event_detector.py /home/pi/Documents/watch_folder --output /home/pi/Documents/processed_images
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# User and Group - change 'pi' to your actual username
User=pi
Group=pi

# Working Directory
WorkingDirectory=/home/pi/Documents/watchmen/event-detect-service

# Environment Variables
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/home/pi/Documents/watchmen/event-detect-service"
Environment="PYTHONUNBUFFERED=1"

# Security settings (adjust as needed)
PrivateNetwork=false
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/pi/Documents
NoNewPrivileges=yes

# Resource limits (optional)
LimitNOFILE=65536
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target